# .deploy/consult-stripe-verify/consult-stripe-verify.service
[Unit]
Description=Consult Stripe signature verifier (edge)
After=network.target
[Service]
Type=simple
WorkingDirectory=/opt/consult-stripe-verify
EnvironmentFile=/opt/consult-stripe-verify/.env
ExecStart=/usr/bin/python3 /opt/consult-stripe-verify/verify.py
Restart=on-failure
User=root
[Install]
WantedBy=multi-user.target
